# plex.conf
[program:plex]
user=plex
environment=HOME="/config"
command=/usr/sbin/start_pms
